# config/includes.chroot/etc/systemd/system/nexus-hardware-profile.service
[Unit]
Description=NexusOS Universal Hardware Profile
DefaultDependencies=no
After=local-fs.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nexus-hardware-profile
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
